# dot_config/systemd/user/gammastep.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/sleep 1

[Install]
WantedBy=
WantedBy=default.target
